# roles/mikeandwan_us/templates/systemd/remote-archive-daily.timer.j2
[Unit]
Description=Run remote archive once a day

[Timer]
OnCalendar=01:30:00

[Install]
WantedBy=timers.target
