# roles/mikeandwan_us/templates/systemd/remote-archive-monthly.service.j2
[Unit]
Description=Long term remote archive job for mikeandwan.us

[Service]
Type=oneshot
ExecStart=/home/svc_www_maw/remote-archive.sh monthly
