# roles/postgres_backup_roles/templates/postgres_backup_roles.timer.j2
[Unit]
Description=Run postgres role backup

[Timer]
OnCalendar={{ pg_backup_roles_ctx.calendar }}

[Install]
WantedBy=timers.target
