# roles/reverse_geocoder/templates/reverse-geocode.timer.j2
[Unit]
Description=Reverse Geocode Timer

[Timer]
OnCalendar={{ reverse_geocode_ctx.calendar }}

[Install]
WantedBy=timers.target
