# roles/reverse_geocoder/templates/reverse-geocode.timer.j2
[Unit]
Description=Reverse Geocode Timer

[Timer]
OnCalendar={{ rev_geo__calendar }}

[Install]
WantedBy=timers.target
