# roles/postgres_backup_database/templates/postgres_backup_database.timer.j2
[Unit]
Description=Run postgres backup for database {{ pg_bu_db__db }}

[Timer]
OnCalendar={{ pg_bu_db__calendar }}

[Install]
WantedBy=timers.target
