# roles/postgres_backup_database/templates/postgres_backup_database.timer.j2
[Unit]
Description=Run postgres backup for database {{ pg_backupdb_ctx.db_name }}

[Timer]
OnCalendar={{ pg_backupdb_ctx.calendar }}

[Install]
WantedBy=timers.target
