# roles/solr_indexer/templates/solr-indexer.timer.j2
[Unit]
Description=Run solr-indexer once a day

[Timer]
OnCalendar={{ solr_indexer_ctx.calendar }}

[Install]
WantedBy=timers.target
