# roles/maw_reverse_geocoder/templates/reverse-geocode.timer.j2
[Unit]
Description=Reverse Geocode Timer

[Timer]
OnCalendar={{ maw_revgeo__calendar }}

[Install]
WantedBy=timers.target
