# roles/maw_solr_indexer/templates/solr-indexer.timer.j2
[Unit]
Description=Run solr-indexer once a day

[Timer]
OnCalendar={{ maw_solridx__calendar }}

[Install]
WantedBy=timers.target
