# roles/postgres_backup_cleanup/templates/postgres_backup_cleanup.timer.j2
[Unit]
Description=Run delete old postgres backups

[Timer]
OnCalendar={{ pg_clean__calendar }}

[Install]
WantedBy=timers.target
