# roles/postgres_backup_cleanup/templates/postgres_backup_cleanup.timer.j2
[Unit]
Description=Run delete old postgres backups

[Timer]
OnCalendar={{ pg_cleanup_ctx.calendar }}

[Install]
WantedBy=timers.target
