# roles/certbot_renew/templates/certbot-renew.timer.j2
[Unit]
Description=certbot_renew timer

[Timer]
OnCalendar={{ certbot_ctx.calendar }}
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
