# roles/postgres_gdrive_archive/templates/postgres_gdrive_archive.timer.j2
[Unit]
Description=Run gdrive archive

[Timer]
OnCalendar={{ postgres_gdrive_archive_ctx.calendar }}

[Install]
WantedBy=timers.target
