# roles/postgres_gdrive_archive/templates/postgres_gdrive_archive.timer.j2
[Unit]
Description=Run gdrive archive

[Timer]
OnCalendar={{ pg_archive__calendar }}

[Install]
WantedBy=timers.target
