# roles/postgres_backup/templates/postgres_backup.timer.j2
[Unit]
Description=Run postgres backup for database {{ db_name }}

[Timer]
OnCalendar={{ calendar }}

[Install]
WantedBy=timers.target
