# roles/postgres_optimize/templates/postgres_optimize.timer.j2
[Unit]
Description=Run vacuumdb with analyze for database {{ pg_optimize_ctx.db_name }}

[Timer]
OnCalendar={{ pg_optimize_ctx.calendar }}

[Install]
WantedBy=timers.target
