# roles/postgres_optimize/templates/postgres_optimization.timer.j2
[Unit]
Description=Run vacuumdb with analyze for database {{ db_name }}

[Timer]
OnCalendar={{ calendar }}

[Install]
WantedBy=timers.target
